# _system/deploy/substrate-check-for-updates.service
# systemd user service for the Substrate update checker (Linux).
#
# Triggered by substrate-check-for-updates.timer every hour.
#
# Placeholders — substituted automatically by `substrate init`:
#   __WORKSPACE__   workspace root path
#   __ENGINE__      engine install path
#   __SUBSTRATE__   substrate CLI binary path

[Unit]
Description=Substrate update checker
After=network-online.target

[Service]
Type=oneshot
ExecStart=__SUBSTRATE__ check-for-updates
WorkingDirectory=__WORKSPACE__
Environment=SUBSTRATE_PATH=__WORKSPACE__
Environment=SUBSTRATE_ENGINE_PATH=__ENGINE__
StandardOutput=append:__WORKSPACE__/_system/logs/check-for-updates-stdout.log
StandardError=append:__WORKSPACE__/_system/logs/check-for-updates-stderr.log
