# _system/deploy/substrate-check-for-updates.timer
# systemd user timer for the Substrate update checker (Linux).
#
# Runs substrate-check-for-updates.service every hour.

[Unit]
Description=Run Substrate update checker every hour

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
